# containers/quadlet/webshot-infra-prodlike.target
[Unit]
Description=Webshot prodlike infra (Quadlet)
Requires=webshot-crawler-network.service
Requires=webshot-postgres.service
Requires=webshot-test-target.service
After=webshot-crawler-network.service
PartOf=webshot-prodlike.target
